# files/local-nvme-autoconf.service
[Unit]
Requires=-.mount

[Service]
Type=oneshot
ExecStart=/opt/system-scripts/bin/local-nvme-autoconf.sh

[Install]
WantedBy=multi-user.target
